# root/etc/s6-overlay/s6-rc.d/conman/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc /usr/sbin/conmand -F -c /config/conman.conf
